# Quantai_Backend/oldbackend/QuantAi-Backend-master/config/supervisor/qi-qcluster.conf
[program:qi-qcluster]
command=/opt/quantai/bin/startqcluster.sh
autostart=true
autorestart=true
startretries=10
user=root
