# Quantai_Backend/oldbackend/QuantAi-Backend-master/config/supervisor/qi-uvicorn.conf
[program:qi-uvicorn]
command=/opt/quantai/bin/startuvicorn.sh
autostart=true
autorestart=true
startretries=10
user=root
